# deploy/com.namche.webhook-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.namche.webhook-relay</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/tashi/sandbox/webhook-relay/dist/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tashi/sandbox/webhook-relay</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/tashi/sandbox/webhook-relay/logs/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/tashi/sandbox/webhook-relay/logs/stderr.log</string>

  <key>UserName</key>
  <string>tashi</string>
</dict>
</plist>
